# roles/fx.datafeed/files/cron.d/fxcm_download_data
# These jobs run a datafeed that keeps the mariadb database data up to date.
#
## Note 1:
### If for some reason the feed stops running and there is a need to catch up and download older data, use fx-shell and run:
### fx-download-fxcm.pl --timeframes=300 --numItems=1440 --verbose

*/5 21-23 * * 0 root sleep 3 && docker start fxcmfeed | grep -v fxcmfeed
*/5 * * * 1-4 root sleep 3 && docker start fxcmfeed | grep -v fxcmfeed
*/5 0-21 * * 5 root sleep 3 && docker start fxcmfeed | grep -v fxcmfeed

* 21-23 * * 0 root docker run -t --rm --link fxdatafeed:fxdatafeed --link fxdataaccounts:fxdataaccounts -v $HOME/fx/cfg:/etc/fxtrader fxtrader/finance-hostedtrader fx-demo-oanda.pl >> /root/account.stats
* * * * 1-4 root docker run -t --rm --link fxdatafeed:fxdatafeed --link fxdataaccounts:fxdataaccounts -v $HOME/fx/cfg:/etc/fxtrader fxtrader/finance-hostedtrader fx-demo-oanda.pl >> /root/account.stats
* 0-21 * * 5 root docker run -t --rm --link fxdatafeed:fxdatafeed --link fxdataaccounts:fxdataaccounts -v $HOME/fx/cfg:/etc/fxtrader fxtrader/finance-hostedtrader fx-demo-oanda.pl >> /root/account.stats
