# roles/fx.datafeed/files/cron.d/oanda_download_data
# These jobs run a datafeed that keeps the mariadb database data up to date.
#
## Note 1:
### If for some reason the feed stops running and there is a need to catch up and download older data, use fx-shell and run:
### fx-download.pl --timeframes=60 --numItems=7200 --verbose --provider=oanda

* 21-23 * * 0 root sleep 3 && docker start oandafeed | grep -v oandafeed
* * * * 1-4 root sleep 3 && docker start oandafeed | grep -v oandafeed
* 0-21 * * 5 root sleep 3 && docker start oandafeed | grep -v oandafeed

10 22 * * * root sleep 3 && docker start oandacleanup | grep -v oandacleanup

15 22 * * * root sleep 3 && docker start oandahistorical | grep -v oandahistorical
